# supervisor-dthb.conf
[program:dthb]
directory=/data/dthb
command=/usr/bin/python3 /data/dthb/run.py
redirect_stderr=true
autostart=true
autorestart=true
startsecs=5
priority=1